# deploy/systemd/fastapi.service
[Unit]
Description=FastAPI service (Uvicorn)
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/app/current/backend
EnvironmentFile=-/opt/app/shared/backend.env
EnvironmentFile=-/opt/app/shared/.env
ExecStart=/opt/app/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
