# deploy/systemd/nextjs.service
[Unit]
Description=Next.js service
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/app/current/frontend
EnvironmentFile=-/opt/app/shared/.env
Environment=PORT=3000
ExecStart=/usr/bin/npm start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
